# ansible/roles/install-containerd/templates/etc/profile.d/containerd.sh
export CONTAINER_RUNTIME_ENDPOINT="unix://{{ containerd.socket_path }}"
