# ansible/roles/install-kubelet/templates/etc/systemd/system/kubelet.service.d/0-cri.conf
[Service]
Environment="KUBELET_EXTRA_ARGS=--container-runtime=docker --runtime-request-timeout=15m --container-runtime-endpoint=unix://{{ docker.socket_path }}"
